# systemd/nora-kiosk.service
[Unit]
Description=Nora Kiosk UI (Chromium)
After=graphical.target nora-backend.service
Wants=nora-backend.service

[Service]
User=nora
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/nora/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
ExecStart=/usr/bin/chromium-browser \
  --kiosk \
  --incognito \
  --noerrdialogs \
  --disable-session-crashed-bubble \
  --disable-infobars \
  --check-for-update-interval=31536000 \
  --overscroll-history-navigation=0 \
  --start-fullscreen \
  --enable-features=OverlayScrollbar \
  http://localhost:8080/
Restart=always
RestartSec=2
Environment=CHROME_LOG_FILE=/home/nora/chromium.log

[Install]
WantedBy=graphical.target
